# systemd/senseuw-site@.timer
[Unit]
Description=Poll %i for updates

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
